# deploy/dashboard-backend.service
[Unit]
Description=Nexus Kiosk Backend
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=KIOSK_USER
WorkingDirectory=INSTALL_DIR
ExecStartPre=/bin/sh -c 'fuser -k 3001/tcp 2>/dev/null; sleep 1; true'
ExecStart=/usr/bin/node INSTALL_DIR/server/dist/index.js
ExecStopPost=/bin/sh -c 'fuser -k 3001/tcp 2>/dev/null; true'
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=10
Environment=NODE_ENV=production
EnvironmentFile=INSTALL_DIR/.env
StandardOutput=append:/var/log/nexus-kiosk/backend.log
StandardError=append:/var/log/nexus-kiosk/backend.log

[Install]
WantedBy=multi-user.target
