# deploy/nexus-kiosk-backup.timer
[Unit]
Description=Nexus Kiosk — Board data backup every 6 hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:15:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
