# deploy/nexus-kiosk-updater.timer
[Unit]
Description=Run VRSI Dashboard auto-update weekly

[Timer]
OnCalendar=Sun 03:30
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
